# crontab/nightjob.sh
#!/bin/bash

export DISPLAY=:0.0

export PATH=$PATH:/usr/bin

cd /home/ali/Workspace/code/sibmusic

python3 main.py

#transmission-gtk
